# foot-server@.socket
[Socket]
ListenStream=%t/foot-%i.sock

[Install]
WantedBy=graphical-session.target
